# config/crontab
*/15 * * * * /usr/bin/curl -s -k http://{{DOMAIN}}/update >/dev/null
